# victoria-metrics/rootfs/etc/services.d/victoria-metrics/run
#!/usr/bin/with-contenv bashio
set -e

RETENTION_PERIOD=$(bashio::config 'retentionPeriod')
ARGS=$(bashio::config 'additionalArguments')

bashio::log.info "Starting Victoria Metrics with retention set to $RETENTION_PERIOD"

if  [ ! -z "$ARGS" ]
then
  bashio::log.info "Starting Victoria Metrics with args set to $ARGS"

fi
exec /victoria-metrics-prod -storageDataPath "/share/victoria-metrics-data" -retentionPeriod $RETENTION_PERIOD $ARGS