# victoria-metrics/rootfs/etc/services.d/victoria-metrics/run
#!/usr/bin/with-contenv bashio
set -e

STORAGE_DATA_PATH=$(bashio::config 'storageDataPath')
RETENTION_PERIOD=$(bashio::config 'retentionPeriod')
USE_PROMETHEUS_CONFIG=$(bashio::config 'usePrometheusConfig')
ARGS=$(bashio::config 'additionalArguments')

bashio::log.info "Starting Victoria Metrics with storage data path set to $STORAGE_DATA_PATH"
bashio::log.info "Starting Victoria Metrics with retention set to $RETENTION_PERIOD"

if  [ ! -z "$ARGS" ]
then
  bashio::log.info "Starting Victoria Metrics with args set to $ARGS"
fi

if  [ "$USE_PROMETHEUS_CONFIG" = true ]
then
  ARGS="${ARGS} -promscrape.config /prometheus.yml"
  bashio::log.info "Starting Victoria Metrics with prometheus scraper"
fi

exec /victoria-metrics-prod -storageDataPath $STORAGE_DATA_PATH -retentionPeriod $RETENTION_PERIOD $ARGS